# conf/sv.redis.ini
[program:redis]
command=/usr/bin/redis-server /etc/redis.conf
user=root
autostart=false
stopsignal=QUIT
